# etc_supervisor_conf.d_mtgoxcachingproxy.conf
[program:mtgoxcachingproxy]
directory=/home/jan/MtGoxCachingProxy
command=java -jar MtGoxCachingProxy.jar
user=jan
redirect_stderr=true
