# etc_supervisor_conf.d_metricsd.conf
[program:metricsd]
directory=/home/jan/metricsd
environment=CONFIG="/home/jan/metricsd/config.json"
command=java -jar metricsd-0.4.2-jar-with-dependencies.jar
user=jan
redirect_stderr=true
