# files/etc/supervisor/conf.d/mtgoxcachingproxy.conf
[program:mtgoxcachingproxy]
directory=/home/jan/MtGoxCachingProxy
command=java -jar MtGoxCachingProxy.jar
user=jan
redirect_stderr=true
startsecs=10
startretries=100000
