# files_graphite/etc/network/if-up.d/iptables
#! /bin/sh

iptables-restore < /etc/iptables/iptables.conf
